# docker/cms/backup.cron
PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin
0 0 * * * cd /home/app && python3 manage.py dbbackup --clean --noinput && python3 manage.py mediabackup --clean --noinput  > /proc/1/fd/1 2>/proc/1/fd/2
